# stow-root/etc/systemd/system/trash-empty.timer
[Unit]
Description=Clean Trash

[Timer]
OnCalendar=daily
Persistent=true
Unit=trash-empty.service

[Install]
WantedBy=timers.target
